# system-files/.local/share/system-files/services/other/etc/runit/sv/timekpr/run
#!/bin/sh
sv check dbus > /dev/null || exit 1
cd /usr/lib/python3/dist-packages/timekpr
exec /usr/bin/timekprd
