# system-files/.local/share/system-files/services/other/etc/runit/sv/skr/run
#!/bin/sh
exec >/dev/null 2>&1
export USER=root
export SUDO_USER=koonix
exec skr
