# system-files/.local/share/system-files/services/other/etc/runit/sv/emby/run
#!/bin/sh
emby_conf=/etc/conf.d/emby-server
export $(grep -v '^#' "$emby_conf" | xargs)
exec chpst -u emby:emby /usr/bin/emby-server
