# system-files/.local/share/system-files/services/useful/etc/runit/sv/windscribe/finish
#!/bin/sh
exec /usr/bin/windscribe stop
