# system-files/.local/share/system-files/services/useful/etc/runit/sv/windscribe/run
#!/bin/sh
sv check NetworkManager > /dev/null || exit 1
/usr/bin/windscribe start
exec pause
