# system-files/.local/share/system-files/services/other/etc/runit/sv/zram/finish
#!/bin/sh
swapoff /dev/zram0
rmmod zram
exec sleep 0.1
