# system-files/.local/share/system-files/services/useful/etc/runit/sv/timesync/run
#!/bin/sh
sv check NetworkManager > /dev/null || exit 1
sleep 5
ntpd -gq
exec pause
