# system-files/.local/share/system-files/services/useful/etc/runit/sv/dictd/run
#!/bin/sh
pkill -x dictd
[ -r /etc/conf.d/dictd ] && . /etc/conf.d/dictd
/usr/bin/dictd $DICTD_ARGS -- $DICTD_EARGS
exec pause
